# ai.planalityc.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.planalityc.autopush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/asans/Documents/табель/planalityc.ai/auto-push.sh</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/asans/Documents/табель/planalityc.ai/.git/COMMIT_EDITMSG</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/asans/Documents/табель/planalityc.ai/.git/auto-push.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/asans/Documents/табель/planalityc.ai/.git/auto-push.log</string>
</dict>
</plist>
